# deploy/vpn-site-api.service.example
[Unit]
Description=RootVPN Website API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpn-bot
EnvironmentFile=/opt/vpn-bot/.env
ExecStart=/opt/vpn-bot/.venv/bin/python /opt/vpn-bot/website_api.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
